# discord_bot/funnel-bot.service
[Unit]
Description=FunnelAgent Discord Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/funnel/FunnelAgent
ExecStart=/home/funnel/FunnelAgent/venv/bin/python -m discord_bot.app.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
